# scripts/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.4.2)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/node /usr/lib/node_modules/openclaw/dist/index.js gateway --port 18789
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group
Environment=HOME=/home/tyranno
Environment=TMPDIR=/tmp
Environment=PATH=/usr/bin:/home/tyranno/.local/bin:/home/tyranno/.npm-global/bin:/home/tyranno/bin:/home/tyranno/.volta/bin:/home/tyranno/.asdf/shims:/home/tyranno/.bun/bin:/home/tyranno/.nvm/current/bin:/home/tyranno/.fnm/current/bin:/home/tyranno/.local/share/pnpm:/usr/local/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.4.2

[Install]
WantedBy=default.target
